# ci/linux/etc/systemd/system/docker-volume-nuke.service
[Unit]
Description=Nuke build cache docker volumes

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/docker-volume-nuke.sh
User=buildkite-agent
