# ci/linux/etc/systemd/system/zfs-prune.timer
[Unit]
Description=Prune ZFS filesystems approaching their quota

[Timer]
OnCalendar=04:00
Persistent=true

[Install]
WantedBy=timers.target
