# ci/linux/etc/systemd/system/zfs-prune.service
[Unit]
Description=Run the ZFS prune script

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/zfs-prune.sh
