# etc/systemd/system/buildkite-agent@.service.d/local.conf
[Service]
EnvironmentFile = /etc/systemd/system/buildkite-agent@.service.d/token.env
