# zockervols/.debian/postinst
#!/usr/bin/env bash
set -eou pipefail

systemctl enable zockervols.socket
systemctl start zockervols.socket
