# zockervols/.debian/prerm
#!/usr/bin/env bash
set -eou pipefail

systemctl stop zockervols.socket
systemctl disable zockervols.socket
