# linux/etc/systemd/system/buildkite-agent@.service.d/local.conf
[Service]
EnvironmentFile=/etc/buildkite-agent/env
