# linux/etc/systemd/system/buildkite-agent@.service.d/local.conf
[Service]
Environment=BUILDKITE_GS_APPLICATION_CREDENTIALS=/etc/gce/cred.json
EnvironmentFile=/etc/buildkite-agent/env
