# ci/etc/systemd/system/docker-volume-prune.timer
[Unit]
Description=Prune build cache docker volumes periodically

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
