# ci/etc/systemd/system/buildkite-agent@.service.d/local.conf
[Service]
EnvironmentFile = /etc/systemd/system/buildkite-agent@.service.d/token.env
Environment=BUILDKITE_AGENT_CONFIG=/etc/buildkite-hooks/buildkite-agent.cfg
