# etc/systemd/system/zockervols.service.d/local.conf
[Service]
Environment="RUST_LOG=info"
User=buildkite-builder
Group=buildkite-builder
