# etc/systemd/system/docker-system-prune.timer
[Unit]
Description=Run docker system prune periodically

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
